# init.d/rb-scanner-request
#!/bin/bash
#
# chkconfig: 2345 60 14
# description: daemon control script
#

/opt/rb/etc/init.d/sv_service rb-scanner-request $*
exit $?

